# systemd/user/battery-notify.service

[Unit]
Description=Battery Low Notification

[Service]
ExecStart=/home/dama/.config/scripts/battery_notify.sh

